# Debugging/libbpf-tools/src/resmon/resmond.service
[Unit]
Description=mlxsw resource monitor

[Service]
Type=simple
ExecStart=/usr/local/bin/resmon start
ExecStop=/usr/local/bin/resmon stop
